# deploy/systemd/homelab-resource-monitor-hub.service
[Unit]
Description=Homelab Resource Monitor telemetry hub
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=homelab-monitor
Group=homelab-monitor
StateDirectory=homelab-resource-monitor
StateDirectoryMode=0750
ExecStart=/opt/homelab-resource-monitor/.venv/bin/python -m hub.app --config /etc/homelab-resource-monitor/hub.json
WorkingDirectory=/opt/homelab-resource-monitor
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
